# conf/supervisord.ini
[program:memcached]
command=/usr/bin/memcached -p 11211 -u memcached -m 256 -c 1024 -t 4
autostart=true
autorestart=true
user=root
redirect_stderr=true

[program:php-fpm]
command = /usr/sbin/php-fpm --nodaemonize
user = root
autostart = true
